# deploy/tailscale-serve.service.d/override.conf
# Drop-in fix for the boot race that left tailscale-serve failed at startup with
# "unexpected state: NoState": the unit ran `tailscale serve` before tailscaled
# had finished coming up. Wait for the backend to report a usable state before
# applying serve rules. Additive — does not modify the base unit.
#
# NOTE: tailscale-serve.service is shared OpenClaw infra (it proxies the gateway
# on :443 and Helm prod on :5174), not a Helm-owned unit — hence only this
# drop-in is vendored here, not the base unit. Install it with:
#
#     mkdir -p ~/.config/systemd/user/tailscale-serve.service.d
#     cp deploy/tailscale-serve.service.d/override.conf \
#        ~/.config/systemd/user/tailscale-serve.service.d/
#     systemctl --user daemon-reload && systemctl --user restart tailscale-serve

[Service]
ExecStartPre=/bin/sh -c 'for i in $(seq 1 30); do /usr/bin/tailscale status >/dev/null 2>&1 && exit 0; sleep 1; done; echo "tailscaled not ready after 30s" >&2; exit 1'
TimeoutStartSec=60
